# coreos/restic.service
[Unit]
Description=Run restic backup

[Service]
Type=oneshot
ExecStart=/var/home/core/backup.sh